# debian/nova-consoleproxy.nova-xenvncproxy.init.in
#!/bin/sh
### BEGIN INIT INFO
# Provides:          nova-xvpvncproxy nova-xenvncproxy
# Required-Start:    $network $local_fs $remote_fs $syslog
# Required-Stop:     $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Nova XVP VNC proxy
# Description:       Nova XVP VNC proxy
### END INIT INFO

# Authors: Thomas Goirand <zigo@debian.org>

DESC="OpenStack XVP VNC proxy"
PROJECT_NAME=nova
NAME=${PROJECT_NAME}-xvpvncproxy

# Exit if the default file not present, or not set to this daemon
if [ -r /etc/default/nova-consoleproxy ] ; then
	. /etc/default/nova-consoleproxy
else
	exit 0
fi
if ! [ "${NOVA_CONSOLE_PROXY_TYPE}" = "xenvnc" ] ; then
	exit 0
fi
